# ansible/roles/secrets/files/sops-secrets.service
[Unit]
Description=Decrypt SOPS secrets at boot
DefaultDependencies=no
# Run before sysusers creates users (they may need passwords)
Before=systemd-sysusers.service sysinit.target
# Require local filesystem to be mounted
After=local-fs.target
# Only run if encrypted secrets exist
ConditionPathExists=/usr/share/secrets/global_secrets.yaml

[Service]
Type=oneshot
ExecStart=/usr/local/bin/decrypt-secrets
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sysinit.target
